# eaa-tool.service
[Unit]
Description=EAA / WCAG Accessibility Check (SCORM)
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/AI_Space/complete/eaa_tool
Environment=SCORM_KIT=/home/ubuntu/AI_Space/complete/autonomous_lab/scorm-kit/bin/scorm-kit.js
Environment=NODE_BIN=/usr/bin/node
ExecStart=/home/ubuntu/.local/bin/gunicorn app:app \
  --workers 2 \
  --worker-class gthread \
  --threads 2 \
  --bind 0.0.0.0:5051 \
  --timeout 120 \
  --keep-alive 5 \
  --max-requests 500 \
  --max-requests-jitter 50 \
  --access-logfile /var/log/eaa-tool-access.log \
  --error-logfile /var/log/eaa-tool-error.log \
  --log-level info
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
